# templates/restic.timer.epp
<%- |
    String $timer,
| -%>
[Timer]
OnCalendar=<%= $timer %>
