# templates/restic.service.epp
<%- |
    $config,
    $group,
    $user,
| -%>
[Service]
User=<%= $user %>
Group=<%= $group %>
Type=oneshot
EnvironmentFile=<%= $config %>
