# templates/restic.service.epp
<%- |
    $commands,
    $config,
    $group,
    $user,
| -%>
[Service]
User=<%= $user %>
Group=<%= $group %>
Type=oneshot
EnvironmentFile=<%= $config %>
<% $commands.each |$command| { -%>
ExecStart=<%= $command %>
<% } -%>
